# systemd/units/eucaconsole-keygen.service
[Unit]
Description=Eucalyptus Management Console key generation
ConditionFileNotEmpty=!/etc/eucaconsole/session-keys.ini

[Service]
ExecStart=/usr/libexec/eucaconsole/eucaconsole-keygen
Type=oneshot
